# home/private_dot_config/systemd/user/jetbrains-installer.timer
[Unit]
Description=Check for latest updates for JetBrains IDEs

[Timer]
OnCalendar=daily
AccuracySec=12h
Persistent=true

[Install]
WantedBy=basic.target